# deploy/timeweb/vrcatalog.service
[Unit]
Description=VR Catalog Docker Compose stack
After=docker.service network.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/var/www/vrcatalog
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
